# agent/acquired.service
[Unit]
Description=response

[Service]
ExecStart=/usr/local/bin/acquired.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
