# agent/acquired.conf
start on (runlevel [345] and started network)
stop on (runlevel [!345] or stopping network)
respawn
exec /usr/local/bin/acquired.py
